# examples/systemd/nexttrace_exporter.service
[Unit]
Description=NextTrace Exporter
Documentation=https://github.com/vinsec/nexttrace_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/nexttrace_exporter \
    --config.file=/etc/nexttrace_exporter/config.yml \
    --web.listen-address=:9101 \
    --log.level=info
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/nexttrace_exporter

# Resource limits
LimitNOFILE=8192
LimitNPROC=512

[Install]
WantedBy=multi-user.target
